# root/etc/services.d/netdata/run
#!/usr/bin/with-contenv bash

exec \
  /usr/sbin/netdata -D  -c /config/netdata/netdata.conf -u abc -s /host -p 19999
